# deploy/tutorbot.service
[Unit]
Description=Tutorbot Flask Application
After=network.target

[Service]
Type=simple
User=stephen
WorkingDirectory=/home/stephen/tutorbot
Environment=PATH=/home/stephen/tutorbot/venv/bin:/usr/local/bin:/usr/bin:/bin
ExecStart=/usr/bin/python3 /home/stephen/tutorbot/main.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target 